# postgres_backup.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/lib/systemd/system/postgres_backup/pg_backup.sh

[Install]
WantedBy=multi-user.target